# config/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:api]
command=/opt/entrypoint.sh
directory=/opt/api
autostart=true
autorestart=true
priority=5
stdout_events_enabled=true
stderr_events_enabled=true
